# files/common/lib/systemd/system/ssh.service.d/override.conf
[Unit]
#
# We want the ssh service to start as early as possible during boot up,
# so we explicitly disable the default dependencies here. This is to
# prevent unrelated services from blocking the ssh service from starting
# (e.g. if an unrelated service failed), as well as enable users to log
# into the system as quickly as possible after a reboot.
#
# We add a single dependency on network.target to ensure that any addresses
# listed as ListenAddress options in /etc/ssh/sshd_config are configured before
# sshd starts because sshd will need to bind explicitly to each.
#
DefaultDependencies=no
After=network.target
