# files/common/lib/systemd/system/zfs-volume-wait.service.d/override.conf
#
# During upgrade verification, the root filesystem will be booted as a
# container using systemd-nspawn. When the zfs-volume-wait service runs
# in this container, it will fail to complete in a timely manner,
# resulting in a failure to start the upgrade container (the logic that
# starts the container will timeout).
#
# The service fails to complete in a timely manner because "/dev" in the
# container is not the same as "/dev" on the host. The ZVOL devices that
# the service is waiting to appear in "/dev", will not appear within the
# container (by design).
#
# Thus, to workaround this problem, we explicitly disable this service
# from running when inside of the container.
#
[Unit]
ConditionVirtualization=!container
