# files/common/lib/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty -t 120 -o '-p -- \\u' --noclear %I $TERM
ExecStartPre=-/usr/bin/delphix-startup-screen
StandardInput=tty
StandardOutput=tty
