# files/common/lib/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty -n -l /usr/bin/delphix-startup-screen --noclear %I $TERM
StandardInput=tty
StandardOutput=tty
