# lib/systemd/system/nfs-server.service.d/override.conf
[Unit]
#
# We want the NFS server to wait to start until after the domain0
# mounts have completed, which occurs via the "zfs-mount" service.
# Otherwise, the NFS exports may come online before the underlying
# datasets are ready to be accessed.
#
After=zfs-mount.service
