# files/common/lib/systemd/system/zfs-import-cache.service.d/override.conf
#
# During upgrade verification, the root filesystem will be booted as a
# container using systemd-nspawn. We don't "sandbox" the container, so
# when zfs-import-cache.service runs, it could potentially try to import
# a pool, which is not desired.
#
# To prevent this behavior, we explicitly disable this service from
# running when inside of the container.
#
[Unit]
ConditionVirtualization=!container
