# files/common/lib/systemd/system/cloud-init.service.d/override.conf
[Unit]
After=local-fs.target
After=nss-lookup.target
ConditionVirtualization=!container
