# files/common/lib/systemd/system/systemd-journald.service.d/override.conf
[Service]
#
# By default, the "systemd-journald" service will use a value of "3min"
# for it's watchdog timer. We've see this timer trigger when running on
# some cloud platforms, when performing certain storage heavy workloads
# on IO limited instances. In these circumstances, the cloud platform
# throttles IO to the point that disk latencies can be on the order of
# seconds, resulting in calls to "fsync" to take minutes.
#
# "systemd-journald" periodically calls "fsync" to flush it's data to
# disk, and as a result of the situation described above, these calls
# can take multiple minutes to complete, and while "systemd-journald" is
# waiting it'll be killed due to this 3 minute timer.
#
# In all the instances we've investigated where this watchdog time limit
# was reached, the system was otherwise adapting to these high disk
# latencies as it should. Thus, we're disabling the watchdog timer here,
# so that "systemd-journald" will also gracefully handle these high disk
# latency events, rather than generating unnecessary core dumps.
#
WatchdogSec=0
