# files/common/etc/update-motd.d/00-delphix
#!/bin/sh

version=$(/usr/bin/get-appliance-version | awk -F'+' '{print $1}')
printf "Welcome to Delphix %s\n" "$version"
