# files/common/etc/update-motd.d/00-delphix
#!/bin/sh

version=$(get-appliance-version --snapshot)
printf "Welcome to Delphix %s\n" "$version"
