# files/common/lib/systemd/system/google-instance-setup.service.d/override.conf
#
# During upgrade verification, the root filesystem will be booted as a
# container using systemd-nspawn. Engines running on GCP have the
# google-instance-setup service enabled by default. When this service
# runs in a container, it fails to activate and cause the start container
# operation to hang, leading to upgrade verification hanging.
# Thus, to workaround this problem, we explicitly disable this service
# from running when inside of the container.
#
[Unit]
ConditionVirtualization=!container
