# debian/preinst
#!/bin/bash -eux
#
# Copyright (c) 2021 by Delphix. All rights reserved.
#

#
# Note that the logic below was copied from delphix-virtualization
# to control whether or not some services should be started.
#
# There's 3 types of services that we account for:
#
#  1. Services that are always enabled.
#  2. Services that are always disabled.
#  3. Services that are disabled by default, but may be
#     dynamically enabled.
#
# delphix-platform currently only needs to deal with services of type 3.
#

case $1 in
install | upgrade)
	#
	# This is a bit of a hack, but we want to ensure these services
	# don't start automatically when the package delivering that
	# service is installed. Unfortunately there isn't a good way to
	# disable services prior to the service being installed; at
	# which point, the service might have already been started.
	#
	# Thus, we rely on this preinst hook to "mask" these services
	# prior to the service being installed. Generally, we'll undo
	# this hack in the postinst hook, by removing the symlink and
	# simply disable the service using the "systemctl" command.
	#

	#
	# Here we account for services of type (3)
	#
	# The systemctl is-enabled will return false if the service is masked,
	# disabled, or if the service doesn't exist yet. By checking this we will
	# only mask the service if it hasn't been installed, is disabled, or is
	# already masked.
	#
	while read -r svc; do
		if ! systemctl is-enabled "$svc"; then
			ln -sf /dev/null "/etc/systemd/system/$svc"
		fi
	done <<-EOF
		ntp.service
	EOF

	;;
esac

exit 0
