# files/aws/lib/systemd/system/systemd-networkd.service.d/20-override.conf
[Service]
ExecStartPre=
ExecStartPre=+/usr/bin/update_netplan_to_mac_address
