# files/aws/lib/systemd/system/systemd-networkd.service.d/aws-override.conf
#
# Copyright 2022 Delphix
#
# Licensed under the Apache License, Version 2.0 (the "License");
# you may not use this file except in compliance with the License.
# You may obtain a copy of the License at
#
#     http://www.apache.org/licenses/LICENSE-2.0
#
# Unless required by applicable law or agreed to in writing, software
# distributed under the License is distributed on an "AS IS" BASIS,
# WITHOUT WARRANTIES OR CONDITIONS OF ANY KIND, either express or implied.
# See the License for the specific language governing permissions and
# limitations under the License.
#

# Changing instance type on AWS could result in the instance booting with an invalid network configuration. See
# DLPX-85209 for more context. This is being handled by switching to MAC based interface names only on AWS from
# 6.0.17.0.
# The fix requires special handling for upgrade scenario. To preserve the network interfaces post reboot into the
# new images, during upgrade workflow on the old image, we create a map file containing a map from existing network
# interface names as keys and their mac address as the value. The map is then referred to by the
# update_netplan_to_mac_address script post reboot to update the network interfaces in the netplan file to their
# corresponding mac address based names. Since this renaming needs to happen before the systemd-networkd executes
# executes its main logic, we run the update_netplan_to_mac_address script as part of PreExecStart directive.

[Service]
ExecStartPre=+/usr/bin/update_netplan_to_mac_address
