# files/aws/lib/systemd/system/systemd-networkd.service.d/override.conf
[Service]
ExecStartPre=
ExecStartPre=+/opt/delphix/server/bin/update_netplan_to_mac_address
