# files/common/lib/systemd/system/nfs-mountd.service.d/override.conf
[Unit]
#
# During upgrade verification, the root filesystem will be booted as a
# container using systemd-nspawn. When this service runs in this
# container, it can cause corruption of the NFS volumes on the client
# (due to us allowing all capabilities to the container).  Thus, to
# workaround this problem, we explicitly disable this service from
# running when inside of the container.
#
ConditionVirtualization=!container
