# etc/systemd/system/nfs-mountd.service.d/override.conf
#
# Copyright 2018 Delphix
#
# Licensed under the Apache License, Version 2.0 (the "License");
# you may not use this file except in compliance with the License.
# You may obtain a copy of the License at
#
#     http://www.apache.org/licenses/LICENSE-2.0
#
# Unless required by applicable law or agreed to in writing, software
# distributed under the License is distributed on an "AS IS" BASIS,
# WITHOUT WARRANTIES OR CONDITIONS OF ANY KIND, either express or implied.
# See the License for the specific language governing permissions and
# limitations under the License.
#

#
# Override the nfs-mountd service execution parameters. By default, the
# service provides the --manage-gids flag so NFS will be aware of group
# membership for users in more than 16 groups. However, this flag has
# the undesirable side-effect of preventing users from accessing files
# owned by a group they are a member of, unless it's their primary
# group. As we do not have a great deal of control of our users' groups
# configuration, we have decided to disable the feature. Since Illumos
# users were already limited to 16 groups for NFS, this is not a change
# in behavior.
#
[Service]
ExecStart=
ExecStart=/usr/sbin/rpc.mountd
