# files/common/lib/systemd/system/rtslib-fb-targetctl.service.d/override.conf
[Unit]
#
# ZFS Volumes (zvols) must be ready when we start this service.
#
Requires=zfs-volumes.target
After=zfs-volumes.target

#
# During upgrade verification, the root filesystem will be booted as a
# container using systemd-nspawn. When the rtslib-fb-targetctl service
# runs in this container, it can cause corruption of the iSCSI volumes
# on the host (due to us allowing all capabilities to the container).
# Thus, to workaround this problem, we explicitly disable this service
# from running when inside of the container.
#
ConditionVirtualization=!container
