# files/common/lib/systemd/system/rtslib-fb-targetctl.service.d/override.conf
[Unit]
#
# This service depends on zvol dev links to be present, which should
# be implemented in DLPX-64724. In the meanwhile, we have this service
# depend on zfs.target, which doesn't guarantee that the zvol links
# will be there, but gives a time buffer after pools have been loaded
# to create the zvol links. This new dependency is especially important
# since local-fs.target, which this service depends on, no longer
# depends on zfs-mount.service.
#
After=zfs.target

#
# During upgrade verification, the root filesystem will be booted as a
# container using systemd-nspawn. When the rtslib-fb-targetctl service
# runs in this container, it can cause corruption of the iSCSI volumes
# on the host (due to us allowing all capabilities to the container).
# Thus, to workaround this problem, we explicitly disable this service
# from running when inside of the container.
#
ConditionVirtualization=!container

[Service]
#
# Although delphix-migration already defines a RequiredBy dependency on
# rtslib-fb-targetctl, that dependency is only honored when delphix-migration
# is enabled. If delphix-migration was manually disabled before migration
# completed, we still want to prevent rtslib-fb-targetctl from running, so
# it is done with this check. Note that "ConditionPathExists" is not a
# valid alternative since the result of the check doesn't prevent
# rtslib-fb-targetctl's dependencies from running.
#
ExecStartPre=/bin/bash -xc "test ! -e /var/delphix/migration/perform-migration"
